# alarm.service
[Unit]
Description=Alarm to wake me up with a song
StartLimitIntervalSec=1

[Service]
Type=simple
Environment=PATH=path
Restart=always
RestartSec=1
ExecStart=homedir/.local/bin/alarm/alarm_service.py

[Install]
WantedBy=multi-user.target
